# app-shack/rootfs/etc/s6-overlay/s6-rc.d/shack/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Shack
# Starts Shack service
# ==============================================================================

# Set log level from add-on configuration (default to INFO)
export LOG_LEVEL=$(bashio::config 'log_level' 'INFO')

bashio::log.info "Starting Shack..."

# Create venv at /data/.venv if it doesn't exist
if [ ! -d "/data/.venv" ]; then
    bashio::log.info "Creating virtual environment at /data/.venv..."
    uv venv /data/.venv
fi

# Initialize uv cache in persistent storage on first run
if [ ! -d "/data/.cache/uv" ]; then
    bashio::log.info "Initializing uv cache in persistent storage..."
    mkdir -p /data/.cache
    cp -r ~/.cache/uv /data/.cache/
fi

# Use persistent cache and set external venv location
export UV_CACHE_DIR=/data/.cache/uv
export UV_PROJECT_ENVIRONMENT=/data/.venv

# Sync base packages using inexact mode to preserve integration packages
# --frozen: use lockfile as-is, don't update
# --inexact: keep packages not in lockfile (integration deps)
# Only sync if needed to minimize output (but more than --quiet)
cd /app || exit 1
if ! uv sync --frozen --no-dev --inexact --check 2>/dev/null; then
    bashio::log.info "Syncing base packages..."
    uv sync --frozen --no-dev --inexact
else
    bashio::log.info "Environment is up to date"
fi

exec uv run python main.py
