# app-greenroom/rootfs/etc/s6-overlay/s6-rc.d/greenroom/finish
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Greenroom
# Handles service exit - S6 will manage the supervision tree
# ==============================================================================
readonly exit_code_service="${1}"
readonly exit_code_signal="${2}"
readonly service="Greenroom"

bashio::log.info \
  "Service ${service} exited with code ${exit_code_service}" \
  "(by signal ${exit_code_signal})"

# Exit with the service's exit code to signal S6
# S6 will handle stopping other services and container shutdown
exit "${exit_code_service}"
