# ha-dasher/rootfs/etc/services.d/dasher/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Dasher
# Starts Dasher
# ==============================================================================

bashio::log.info "Starting Dasher..."

exec /app/dasher.py
