# app-shack/rootfs/etc/services.d/shack/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Shack
# Starts Shack service
# ==============================================================================

bashio::log.info "Starting Shack..."

# Create venv at /data/.venv if it doesn't exist
if [ ! -d "/data/.venv" ]; then
    bashio::log.info "Creating virtual environment at /data/.venv..."
    uv venv /data/.venv
fi

# Don't use uv sync to maintain packages installed for integrations
bashio::log.info "Installing base packages..."
uv pip install -r /app/pyproject.toml --python /data/.venv/bin/python

exec /data/.venv/bin/python /app/main.py
