# app-shack/rootfs/etc/services.d/shack/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Shack
# Starts Shack service
# ==============================================================================

bashio::log.info "Starting Shack..."

uv run --directory /app main.py
